# conf/etc/service/brubeck/run
#!/bin/bash

/bin/generate_config.sh

exec /bin/brubeck --config=/root/config.json
