# conf/etc/service/carbonapi/run
#!/bin/bash

sleep 5

exec /usr/bin/carbonapi -config /etc/carbonapi/carbonapi.yaml
