# conf/etc/service/go-carbon/run
#!/bin/bash

exec /usr/bin/go-carbon -config=/etc/go-carbon/go-carbon.conf -daemon=false 2>&1 >> /var/log/go-carbon.log
